# systemd/system/kill-printd.service
Description= Kill fprintd before sleep
Before=sleep.target

[Service]
ExecStart=/bin/killall fprintd

[Install]
WantedBy=sleep.target
