# systemd/system/nftables-sshd-exception.service
[Unit]
Description=adds a vpn network exception for sshd

[Service]
Type=oneshot
ExecStart=/usr/bin/nft -f /home/arts/aurGit/dotFiles/scripts/myScripts/sshd-exclude-traffic.rules

[Install]
WantedBy=network-online.target
